# systemd/getty@.service.d/disable-cursor-blink.conf
[Service]
ExecStartPost=/bin/sh -c 'echo 0 | tee /sys/class/graphics/fbcon/cursor_blink'
